# recipes-nonlinear/nonlinear-rootfs/files/etc/udev/rules.d/100-gether.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="1d6b", ATTR{idProduct}=="0002", ACTION=="add", RUN+="/bin/systemctl start usbnet.service"
